# .config/systemd/user/feh-randomizing-wallpaper.timer
[Unit]
Description=Realoading wallpaper and selecting a random image

[Timer]
OnStartupSec=2s
OnUnitInactiveSec=15m

[Install]
WantedBy=timers.target
